# supervisord.conf
[supervisord]
nodaemon=true

[program:npm]
command=npm start
directory=/data
user=nobody
autostart=true
autorestart=true
redirect_stderr=true
